# build/supervisor.conf
[program:django]
process_name=%(program_name)s
numprocs=1
directory=/root/src/
command=gunicorn aor.wsgi:application
        --bind 0.0.0.0:80
        --workers 2
        --log-config "/root/gunicorn_logging.conf"
stderr_logfile=/root/logs/aor_sv_guni-error.log
stdout_logfile=/root/logs/aor_sv_guni-out.log
autostart=true
autorestart=true
