# build/supervisor.conf
[program:django]
process_name=%(program_name)s
numprocs=1
directory=/root/src/
command=gunicorn aor.wsgi:application
        -b 0.0.0.0:80
        -w 2
        --access-logfile "/root/logs/aor.log"
        --error-logfile "/root/logs/aor-error.log"
stderr_logfile=/root/logs/aor_sv_guni-error.log
stdout_logfile=/root/logs/aor_sv_guni-out.log
autostart=true
autorestart=true
